# conf/magento.cron
*/5 * * * * php /home/magento/files/html/bin/magento cron:run >> /home/magento/files/html/var/log/magento.cron.log
